# files/uwsgi.service
[Unit]
Description=Lisp UWSGI service
After=

[Service]
ExecStart=/lisp/start-uwsgi.sh

[Install]
WantedBy=multi-user.target